# scripts/launchd/com.pocketpod.server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.pocketpod.server</string>

  <key>ProgramArguments</key>
  <array>
    <string>__POCKET_POD_HOME__/.venv/bin/python</string>
    <string>__POCKET_POD_HOME__/server.py</string>
  </array>

  <key>WorkingDirectory</key>
  <string>__POCKET_POD_HOME__</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>POCKET_POD_BASE_URL</key>
    <string>__POCKET_POD_BASE_URL__</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>__POCKET_POD_HOME__/data/logs/server.out.log</string>

  <key>StandardErrorPath</key>
  <string>__POCKET_POD_HOME__/data/logs/server.err.log</string>
</dict>
</plist>
